# scripts/com.stockanalyzer.telegram.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stockanalyzer.telegram</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/puneetgrover/stock-analyzer/venv310/bin/python</string>
        <string>/Users/puneetgrover/stock-analyzer/data_collectors/telegram_bot.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/puneetgrover/stock-analyzer</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DATABASE_URL</key>
        <string>postgresql://puneetgrover@localhost/stock_analyzer</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/puneetgrover/stock-analyzer/logs/telegram_bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/puneetgrover/stock-analyzer/logs/telegram_bot.log</string>
</dict>
</plist>
